# docker/php/crontab
* * * * * www-data /bin/sh /var/www/html/cron.sh cron.php -m=default
* * * * * www-data /bin/sh /var/www/html/cron.sh cron.php -m=always
